# scripts/systemd/rove-snapshot.timer
# Daily backup-snapshot pass for rove-snapshot.service.
#
# Fires once a day at 02:00 local with up to 1h of jitter so a fleet
# of nodes doesn't thunder-herd S3 PUTs. `Persistent=true` catches up
# after a gap (missed timer runs fire on next boot / on next user-
# manager start) — useful for laptops or hosts that sleep through the
# scheduled hour.
#
# Each capture is a leader-only event (the .service file's wrapper
# probes /_system/leader; followers exit 0 silently). So on a 3-node
# cluster all three timers fire but only one node uploads — minimal
# coordination needed.
#
# Adjust the cadence as needed:
#   OnCalendar=hourly      better RPO; ~24× the S3 PUT count
#   OnCalendar=*-*-* 02,14:00:00   twice-daily
#
# Inspect:
#   systemctl --user list-timers rove-snapshot.timer
#   journalctl --user -u rove-snapshot.service -n 50
#
# Install: see scripts/systemd/rove-snapshot.service for the full
# recipe.

[Unit]
Description=Daily Loop46 snapshot capture

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=1h
Persistent=true
Unit=rove-snapshot.service

[Install]
WantedBy=timers.target
