# scripts/systemd/rove-log-server.service
# User-scoped service unit for log-server-standalone.
#
# Companion units:
#   rove-files-server.service  — files-server-standalone.
#   rove-loop46.service        — the worker. Depends on this unit.
#   rove-sse-server.service    — sse-server-standalone.
#
# Full deployment recipe: docs/deployment.md.
#
# ── Install (per host, as the deploy user, no sudo) ──────────────────
# 1. Place the binary on PATH:
#      install -D -m 0755 zig-out/bin/log-server-standalone \
#          ~/.local/bin/log-server-standalone
#
# 2. Place the unit file:
#      install -D -m 0644 scripts/systemd/rove-log-server.service \
#          ~/.config/systemd/user/rove-log-server.service
#
# 3. Drop the env file. log-server reads its log batch store from
#    the same S3 backend as the worker; LOG_S3_KEY_PREFIX scopes the
#    batches to a named subdirectory.
#      cat > ~/.config/rove/log-server.env <<EOF
#      LOG_LISTEN=127.0.0.1:8445
#      DATA_DIR=%h/.rove/data
#      TLS_CERT=%h/.rove/tls/cert.pem
#      TLS_KEY=%h/.rove/tls/key.pem
#      ADMIN_ORIGIN=https://app.loop46.me
#      LOG_S3_KEY_PREFIX=logs/
#      LOG_POLL_INTERVAL_MS=500
#      EOF
#      chmod 0600 ~/.config/rove/log-server.env
#
#    common.env (shared with the other services) must also exist —
#    see docs/deployment.md or rove-files-server.service for its shape.
#
# 4. Enable + start:
#      systemctl --user daemon-reload
#      systemctl --user enable --now rove-log-server.service
#      journalctl --user -u rove-log-server.service -f

[Unit]
Description=Loop46 log-server (per-tenant request log indexer + query)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/logs-plan.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/rove/common.env
EnvironmentFile=%h/.config/rove/log-server.env
ExecStart=%h/.local/bin/log-server-standalone \
    --data-dir ${DATA_DIR} \
    --listen ${LOG_LISTEN} \
    --tls-cert ${TLS_CERT} \
    --tls-key ${TLS_KEY} \
    --cors-origin ${ADMIN_ORIGIN} \
    --poll-interval-ms ${LOG_POLL_INTERVAL_MS}
Restart=on-failure
RestartSec=2

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.rove
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
