# scripts/systemd/rove-sse-server.service
# User-scoped service unit for sse-server-standalone.
#
# Companion units:
#   rove-files-server.service  — files-server-standalone.
#   rove-log-server.service    — log-server-standalone.
#   rove-loop46.service        — the worker. Depends on this unit.
#
# Full deployment recipe: docs/deployment.md.
#
# ── Install (per host, as the deploy user, no sudo) ──────────────────
# 1. Place the binary on PATH:
#      install -D -m 0755 zig-out/bin/sse-server-standalone \
#          ~/.local/bin/sse-server-standalone
#
# 2. Place the unit file:
#      install -D -m 0644 scripts/systemd/rove-sse-server.service \
#          ~/.config/systemd/user/rove-sse-server.service
#
# 3. Drop the env file. sse-server doesn't talk to S3 — it only
#    accepts emit POSTs from the worker (bearer-checked against
#    SSE_INTERNAL_TOKEN) and fans out to subscribed clients over
#    TLS-terminated h2. The clients authenticate with JWTs minted by
#    the worker, verified here against LOOP46_SERVICES_JWT_SECRET.
#      cat > ~/.config/rove/sse-server.env <<EOF
#      SSE_LISTEN=127.0.0.1:8446
#      TLS_CERT=%h/.rove/tls/cert.pem
#      TLS_KEY=%h/.rove/tls/key.pem
#      SSE_INTERNAL_TOKEN=...       # any opaque bearer; must match the
#                                   # worker's SSE_INTERNAL_TOKEN env
#      EOF
#      chmod 0600 ~/.config/rove/sse-server.env
#
#    common.env (shared) must define LOOP46_SERVICES_JWT_SECRET.
#
# 4. Enable + start:
#      systemctl --user daemon-reload
#      systemctl --user enable --now rove-sse-server.service
#      journalctl --user -u rove-sse-server.service -f

[Unit]
Description=Loop46 sse-server (HTTP/2 SSE fan-out from worker emits)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/sse-plan.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/rove/common.env
EnvironmentFile=%h/.config/rove/sse-server.env
ExecStart=%h/.local/bin/sse-server-standalone \
    --listen ${SSE_LISTEN} \
    --tls-cert ${TLS_CERT} \
    --tls-key ${TLS_KEY}
Restart=on-failure
RestartSec=2

ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
