# scripts/systemd/rove-snapshot.service
# User-scoped oneshot for periodic backup snapshots.
#
# Companion timer: rove-snapshot.timer (daily fire).
#
# What it does: invokes `scripts/rove-snapshot-pass.sh`. That wrapper
# probes `/_system/leader` on the local worker; only the leader runs
# the actual `loop46 snapshot` capture (uploads VACUUM-INTO copies
# of every tenant's app.db + __root__.db to S3 under
# `cluster/snapshots/{snap_id}/`). Followers exit 0 silently.
#
# Full deployment recipe: docs/deployment.md "Backups" section.
#
# ── Install (per host, as the deploy user, no sudo) ──────────────────
# Same install as the four binaries — see rove-loop46.service for the
# general pattern. After binaries are placed:
#
# 1. Place the wrapper script:
#      install -D -m 0755 scripts/rove-snapshot-pass.sh \
#          ~/.local/bin/rove-snapshot-pass.sh
#
# 2. Place the unit + timer:
#      install -D -m 0644 scripts/systemd/rove-snapshot.service \
#          ~/.config/systemd/user/rove-snapshot.service
#      install -D -m 0644 scripts/systemd/rove-snapshot.timer \
#          ~/.config/systemd/user/rove-snapshot.timer
#
# 3. Enable + start the timer:
#      systemctl --user daemon-reload
#      systemctl --user enable --now rove-snapshot.timer
#      systemctl --user list-timers rove-snapshot.timer
#
# 4. Force a first capture to confirm wiring:
#      systemctl --user start rove-snapshot.service
#      journalctl --user -u rove-snapshot.service -n 30
#
# Retention: configure an S3 lifecycle policy on the bucket to delete
# `cluster/snapshots/*` keys older than the retention window. See
# docs/deployment.md for the JSON.

[Unit]
Description=Loop46 periodic snapshot capture (leader-only)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/deployment.md
After=network-online.target rove-loop46.service
# Don't fail if the worker isn't up — the wrapper handles that case
# (curl returns 000, we skip). Wants= ties the lifecycle without a
# hard ordering dependency.
Wants=network-online.target rove-loop46.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/rove/common.env
EnvironmentFile=%h/.config/rove/loop46.env
ExecStart=%h/.local/bin/rove-snapshot-pass.sh

# Hardening.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.rove
PrivateTmp=true
NoNewPrivileges=true
