# scripts/systemd/rove-loop46.service
# User-scoped service unit for the loop46 worker.
#
# Companion units (start order: standalones → worker):
#   rove-files-server.service  — files-server-standalone.
#   rove-log-server.service    — log-server-standalone.
#   rove-sse-server.service    — sse-server-standalone.
#
# Full deployment recipe: docs/deployment.md.
#
# ── Install (per host, as the deploy user, no sudo) ──────────────────
# 1. Place the binary on PATH:
#      install -D -m 0755 zig-out/bin/loop46 ~/.local/bin/loop46
#
# 2. Bind privileged port (one of):
#      sudo sysctl -w net.ipv4.ip_unprivileged_port_start=443
#      # persist: /etc/sysctl.d/99-rove.conf
#    OR per-binary capability (re-apply after each `zig build install`):
#      sudo setcap cap_net_bind_service=+ep ~/.local/bin/loop46
#
# 3. Place the unit file:
#      install -D -m 0644 scripts/systemd/rove-loop46.service \
#          ~/.config/systemd/user/rove-loop46.service
#
# 4. Drop the env file. Single-node defaults below; for a multi-node
#    cluster set NODE_ID per host and list all peers in PEERS.
#      cat > ~/.config/rove/loop46.env <<EOF
#      NODE_ID=0
#      PEERS=127.0.0.1:40146
#      RAFT_LISTEN=127.0.0.1:40146
#      HTTP_LISTEN=0.0.0.0:443
#      DATA_DIR=%h/.rove/data
#      TLS_CERT=%h/.rove/tls/cert.pem
#      TLS_KEY=%h/.rove/tls/key.pem
#      PUBLIC_SUFFIX=rewindjs.app          (customer wildcard)
#      SYSTEM_SUFFIX=rewindjs.com          (system surfaces; MUST
#                                            differ from PUBLIC_SUFFIX)
#      ADMIN_API_DOMAIN=app.rewindjs.com   (admin = system surface)
#      ADMIN_ORIGIN=https://app.rewindjs.com
#      FILES_PUBLIC_BASE=https://files.loop46.me:8444
#      LOG_PUBLIC_BASE=https://logs.loop46.me:8445
#      SSE_PUBLIC_BASE=https://sse.loop46.me:8446
#      WORKERS=0                          # 0 = nCPU - 1
#      LOOP46_ROOT_TOKEN=...              # 64 hex chars; bearer auth
#      SSE_INTERNAL_TOKEN=...             # must match sse-server's env
#      EOF
#      chmod 0600 ~/.config/rove/loop46.env
#
#    common.env (shared with standalones) must also be present —
#    see docs/deployment.md or rove-files-server.service for its shape.
#
# 5. Enable + start (after the three standalones are running):
#      systemctl --user daemon-reload
#      systemctl --user enable --now rove-loop46.service
#      journalctl --user -u rove-loop46.service -f

[Unit]
Description=Loop46 worker (HTTP/2 + raft + per-tenant JS dispatch)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/PLAN.md
After=network-online.target rove-files-server.service rove-log-server.service rove-sse-server.service
Wants=network-online.target rove-files-server.service rove-log-server.service rove-sse-server.service

[Service]
Type=simple
EnvironmentFile=%h/.config/rove/common.env
EnvironmentFile=%h/.config/rove/loop46.env
ExecStart=%h/.local/bin/loop46 worker \
    --node-id ${NODE_ID} \
    --peers ${PEERS} \
    --listen ${RAFT_LISTEN} \
    --http ${HTTP_LISTEN} \
    --data-dir ${DATA_DIR} \
    --tls-cert ${TLS_CERT} \
    --tls-key ${TLS_KEY} \
    --public-suffix ${PUBLIC_SUFFIX} \
    --system-suffix ${SYSTEM_SUFFIX} \
    --admin-api-domain ${ADMIN_API_DOMAIN} \
    --admin-origin ${ADMIN_ORIGIN} \
    --files-public-base ${FILES_PUBLIC_BASE} \
    --log-public-base ${LOG_PUBLIC_BASE} \
    --sse-public-base ${SSE_PUBLIC_BASE} \
    --workers ${WORKERS}
# Restart=always (not on-failure) is load-bearing for far-behind
# follower auto-recovery. When this node falls behind the leader's
# raft compaction floor, the leader sends `snap_fetch_offer`; the
# follower's fetcher thread downloads the snapshot bundle, stages it
# under `{data_dir}/.snap-in-{snap_id}/`, then calls
# `std.process.exit(0)` (see src/loop46/snapshot.zig). systemd
# restarts the process; the next boot's `installStagedSnapshotIfPresent`
# atomic-renames the staged files into place and `raft_load_snapshot`
# picks up at the leader's floor. With Restart=on-failure the clean
# exit wouldn't fire a restart and operators would be on the hook.
Restart=always
RestartSec=2

# A worker holding raft state shouldn't be terminated mid-commit.
# 30s is enough for an in-flight propose round to finish before
# systemd SIGKILLs.
TimeoutStopSec=30

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.rove
PrivateTmp=true
NoNewPrivileges=true
# cap_net_bind_service must be granted on the binary itself (step 2
# above) or via sysctl. AmbientCapabilities here doesn't help because
# user units can't add caps not held by the user manager.

[Install]
WantedBy=default.target
