# scripts/systemd/v2/rewind-front.service
# User-scoped service unit for the V2 stateless edge (rewind-front).
# See rewind-cp.service for the install recipe + the 3-process picture.
# Full plan: docs/v2-production-deploy-plan.md.
#
# The front takes ONE positional arg: <tls_port> (443). It also opens a
# plaintext :80 listener (ACME HTTP-01 + HTTP→HTTPS redirect) via
# REWIND_HTTP_PORT. Both are privileged ports, so on each host do ONE of:
#
#   (recommended, survives redeploys — one-time per host)
#     echo 'net.ipv4.ip_unprivileged_port_start=80' | sudo tee /etc/sysctl.d/99-rove.conf
#     sudo sysctl --system
#
#   (alternative — must be RE-APPLIED after every binary replacement;
#    scripts/deploy.sh does this when ROVE_SETCAP=1)
#     sudo setcap cap_net_bind_service=+ep ~/.local/bin/rewind-front
#
# AmbientCapabilities= does NOT help here: user units can't add caps the
# user manager doesn't already hold.
#
# TLS: REWIND_TLS_CERT/KEY is the default platform-wildcard cert (set
# below via %h); per-host custom-domain certs sync from the CP at
# runtime (CertSync). The front holds NO raft state — restart anytime.

[Unit]
Description=rewind-front (V2 stateless edge: public TLS termination + routing)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/v2-production-deploy-plan.md
After=network-online.target rewind-cp.service rewind-worker.service
Wants=network-online.target rewind-cp.service rewind-worker.service

[Service]
Type=simple
EnvironmentFile=%h/.config/rove/common.env
# Default-cert paths need %h, so they're set here (not in common.env,
# which can't expand specifiers). Drop the cert/key at these paths.
Environment=REWIND_TLS_CERT=%h/.rove/tls/platform.crt
Environment=REWIND_TLS_KEY=%h/.rove/tls/platform.key
ExecStart=%h/.local/bin/rewind-front 443

# Stateless — fast restart, no quorum to rejoin.
Restart=always
RestartSec=1
TimeoutStopSec=10

# TLS termination does crypto on the request path; give it a middling
# share between the worker (bulk) and the idle CP.
CPUWeight=3000
MemoryHigh=4G
LimitNOFILE=524288

# Hardening. Reads the TLS cert under %h/.rove; writes nothing to disk.
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=%h/.rove
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
