# scripts/systemd/v2/rewind-worker.service
# User-scoped service unit for the V2 data-plane worker (rewind-worker).
# See rewind-cp.service for the install recipe + the 3-process picture.
# Full plan: docs/v2-production-deploy-plan.md.
#
# The worker takes TWO positional args: <data_dir> <http_port>.
# Everything else is env (common.env + node.env). It listens PLAINTEXT
# h2c on :8443 — the front terminates public TLS and dials it privately,
# so :8443 (and the raft peer port in REWIND_PEERS, :8501) MUST be
# firewalled to the 3 nodes only. There is no app-layer auth on these.

[Unit]
Description=rewind worker (V2 data plane: multi-raft + per-tenant JS dispatch)
Documentation=https://github.com/anarchodev/rove/blob/main/docs/v2-production-deploy-plan.md
After=network-online.target rewind-cp.service
Wants=network-online.target rewind-cp.service

[Service]
Type=simple
EnvironmentFile=%h/.config/rove/common.env
EnvironmentFile=%h/.config/rove/node.env
ExecStartPre=/usr/bin/install -d -m0700 %h/.rove/data/worker
ExecStart=%h/.local/bin/rewind-worker %h/.rove/data/worker 8443

# Restart=always is load-bearing for boot-time tenant-group recovery: a
# restarted node re-joins every tenant raft group it hosts and catches
# up from quorum (commit dd836b4). With on-failure a clean exit wouldn't
# trigger the rejoin restart.
Restart=always
RestartSec=2
# Holds per-tenant raft state; allow an in-flight propose round to
# finish before SIGKILL.
TimeoutStopSec=30

# The worker is the heavy process — give it the bulk of the CPU share.
CPUWeight=10000
MemoryHigh=48G

# Per-worker io_uring rings + registered fds + per-tenant LMDB maps blow
# through the 1024 default soft limit at boot (UserFdQuotaExceeded).
LimitNOFILE=524288

# Hardening.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.rove
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
